# Debian Build/DEBIAN/postinst
#!/bin/bash
set -euo pipefail

if grep -q "boot=live" /proc/cmdline || [ -f /run/live/medium/live.cfg ]; then
    echo "live environment detected — skipping user install"
    exit 0
fi
chmod +x /usr/share/once/once.sh
chmod +x /usr/bin/once
chown root:root /usr/bin/once

for user_dir in /home/*; do
    [ -d "$user_dir" ] || continue
    username=$(basename "$user_dir")
    autostart_dir="$user_dir/.config/autostart"
    mkdir -p "$autostart_dir"
    cp /etc/skel/.config/autostart/once.desktop "$autostart_dir/"
    chown -R "$username:$username" "$autostart_dir"
done

systemctl enable once
